# templates/borgmatic.timer.j2
[Unit]
Description=borgmatic backup timer for {{ item.name }}

[Timer]
OnCalendar={{ item.schedule.backup.oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
